# core-systems/offensive-security-core/reporting/mitre_mapping/tactics_to_techniques.yaml
# TeslaAI Genesis - Промышленная MITRE ATT&CK Mapping
# Используется для AI-анализаторов, генераторов отчётов и построения Kill-Chain цепочек

meta:
  version: "v4.7-industrial"
  source_framework: "MITRE ATT&CK v14 Enterprise"
  last_updated: "2025-07-30"
  generated_by: "TeslaAI Threat Intelligence Fabric"
  validated_by: 
    - "RedOps Consilium (20 agents)"
    - "MetaGenerals Group (3 signatures)"

mappings:
  Initial Access:
    - id: T1190
      name: Exploit Public-Facing Application
      detection_tags: [web-exploit, unauthenticated-access, cve_trigger]
    - id: T1133
      name: External Remote Services
      detection_tags: [rdp-access, ssh-login, vpn-abuse]
    - id: T1078
      name: Valid Accounts
      detection_tags: [creds-abuse, lateral-creds]

  Execution:
    - id: T1059.001
      name: Command and Scripting Interpreter: PowerShell
      detection_tags: [powershell, obfuscated-script]
    - id: T1053.005
      name: Scheduled Task/Job: Scheduled Task
      detection_tags: [task-scheduler, persist-method]

  Persistence:
    - id: T1547.001
      name: Registry Run Keys / Startup Folder
      detection_tags: [autorun, registry-persistence]
    - id: T1136.001
      name: Create Account: Local Account
      detection_tags: [account-creation, privilege-persistence]

  Privilege Escalation:
    - id: T1068
      name: Exploitation for Privilege Escalation
      detection_tags: [exploit-sudo, kernel-eop]
    - id: T1134.001
      name: Access Token Manipulation: Token Impersonation
      detection_tags: [impersonation, sid-theft]

  Defense Evasion:
    - id: T1562.001
      name: Disable or Modify Tools: Security Software
      detection_tags: [av-kill, edr-bypass]
    - id: T1027
      name: Obfuscated Files or Information
      detection_tags: [packers, base64, xor-encoding]

  Credential Access:
    - id: T1003.001
      name: OS Credential Dumping: LSASS Memory
      detection_tags: [lsass-dump, mimikatz]
    - id: T1555.003
      name: Credentials from Web Browsers
      detection_tags: [browser-creds, token-theft]

  Discovery:
    - id: T1082
      name: System Information Discovery
      detection_tags: [os-enum, hostname-grab]
    - id: T1046
      name: Network Service Scanning
      detection_tags: [nmap, portscan]

  Lateral Movement:
    - id: T1021.002
      name: Remote Services: SMB/Windows Admin Shares
      detection_tags: [smb-share, admin-share-access]
    - id: T1075
      name: Pass the Hash
      detection_tags: [pth, ntlm-replay]

  Collection:
    - id: T1113
      name: Screen Capture
      detection_tags: [screensnap, user-monitoring]
    - id: T1056.001
      name: Input Capture: Keylogging
      detection_tags: [keylog, input-hook]

  Exfiltration:
    - id: T1041
      name: Exfiltration Over C2 Channel
      detection_tags: [exfil-c2, tcp-stealth]
    - id: T1567.002
      name: Exfiltration to Cloud Storage
      detection_tags: [dropbox, s3-upload]

  Command and Control:
    - id: T1071.001
      name: Application Layer Protocol: Web Protocols
      detection_tags: [http-c2, tls-c2]
    - id: T1095
      name: Non-Application Layer Protocol
      detection_tags: [icmp-c2, dns-tunnel]

visualization:
  color_by_tactic:
    Initial Access: "#ff5252"
    Execution: "#ff9800"
    Persistence: "#ffd54f"
    Privilege Escalation: "#fbc02d"
    Defense Evasion: "#cddc39"
    Credential Access: "#8bc34a"
    Discovery: "#4caf50"
    Lateral Movement: "#03a9f4"
    Collection: "#00bcd4"
    Exfiltration: "#9c27b0"
    Command and Control: "#673ab7"

integration:
  enable_cvss_linking: true
  enable_heatmap_overlay: true
  include_external_ids: true
  validate_ids: true
  validate_tags: true

